# server/ansible/roles/go/files/hippoboard_api.service.j2
[Unit]
Description=Go REST API Server

[Service]
ExecStart=/var/www/{{ http_host }}/api/api
WorkingDirectory=/var/www/{{ http_host }}/api
Restart=always

[Install]
WantedBy=multi-user.target
